# deploy/supervisor/horizon.conf
; SocialPulse — Laravel Horizon
; Copiar a /etc/supervisor/conf.d/socialpulse-horizon.conf
; sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start socialpulse-horizon

[program:socialpulse-horizon]
process_name=%(program_name)s
command=php /var/www/socialpulse/artisan horizon
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/socialpulse/storage/logs/horizon.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=3600
stopasgroup=true
killasgroup=true
